# docker/supervisord.conf
[supervisord]
nodaemon=true
[program:nlu]
command=/home/homo/homo/nlu.sh

[program:homo-webview]
command=/home/homo/homo/homo-webview